# bridge/launchd/com.weirdapps.telegram-claude-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.weirdapps.telegram-claude-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/plessas/SourceCode/telegram-bot/bridge/launchd/run.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/plessas/SourceCode/telegram-bot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/plessas</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/plessas/Library/Logs/telegram-claude-bridge.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/plessas/Library/Logs/telegram-claude-bridge.err.log</string>
</dict>
</plist>
